# github-bot.conf
description "github-bot"

# used to be: start on startup
# until we found some mounts weren't ready yet while booting:
start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    cd /location/of/
    exec /usr/bin/node server.js >> /var/log/github-bot.log 2>&1
end script
